# local_camera_slave.service
[Unit]
Description=GERTIE Qt Local Camera Service (rep8 on control1)
After=network.target

[Service]
Type=simple
User=andrc1
Group=andrc1
WorkingDirectory=/home/andrc1/camera_system_qt_conversion
ExecStart=/usr/bin/python3 /home/andrc1/camera_system_qt_conversion/local_camera_slave.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=PYTHONPATH=/home/andrc1/camera_system_qt_conversion
Environment=DISPLAY=:0

[Install]
WantedBy=multi-user.target
